# build/hostname/hostname-changer.service
[Unit]
Description=Service to change the hostname
Wants=network-pre.target
Before=network-pre.target

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/usr/local/share/scripts/
ExecStart=/bin/bash -c "/usr/local/share/scripts/hostname-changer.sh"

[Install]
WantedBy=network.target